# .config/autostart/autostart.desktop
[Desktop Entry]
Type=Application
Version=1.0
Name=Autostart
Comment=Run autostart.sh on login
Exec=autostart.sh
